# build/rootfs/etc/supervisor/conf.d/nginx.conf
[program:nginx]
priority=10
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
logfile=/home/container/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
childlogdir=/home/container/log ; ('AUTO' child log dir, default $TEMP)

